# debian/clearwater-socket-factory.upstart
# clearwater-socket-factory
#
# This process runs in the default namespace, and enables processes running in
# different namespaces to obtain sockets in the default namespace.

description "Clearwater socket factory server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
  . /etc/clearwater/config
  /usr/share/clearwater/bin/clearwater_socket_factory "$sas_server"
end script
