# debian/clearwater-socket-factory.upstart
# clearwater-socket-factory
#
# This process runs as root in the default namespace, and enables processes
# running in different namespaces to obtain sockets in other namespaces.

description "Clearwater socket factory server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
  # Sleep for 5 seconds, to avoid respawning socket-factory too quickly
  sleep 5

  /usr/share/clearwater/bin/clearwater-socket-factory-wrapper
end script
